# luci-app-haproxy-tcp/root/etc/uci-defaults/z99-haproxy
#!/bin/sh

/etc/init.d/haproxy disable
/etc/init.d/haproxy stop

rm -f /tmp/luci-indexcache
killall -HUP rpcd 2>/dev/null
exit 0
