# luci-app-argon-config/root/etc/uci-defaults/luci-argon-config
#!/bin/sh

sed -i 's/cbi.submit\"] = true/cbi.submit\"] = \"1\"/g' /usr/lib/lua/luci/dispatcher.lua

rm -f /tmp/luci-indexcache
killall -HUP rpcd 2>/dev/null
exit 0
